# packages/rotator/crontab
0,15,30,25 * * * * node /usr/src/app/packages/rotator/dist/index.js
